# gitlab-ci/iopsys-supervisord.conf
[program:ubusd]
priority=1
autorestart=false
startretries=0
command=/bin/bash -c "/usr/sbin/ubusd"

[program:wifimngr]
priority=2
autorestart=false
startretries=0
command=/bin/bash -c "/usr/sbin/wifimngr"

[program:ieee1905d]
priority=3
autorestart=false
startretries=0
command=/bin/bash -c "/usr/sbin/ieee1905d"

[program:mapagent]
priority=4
autorestart=false
startretries=0
command=/bin/bash -c "/usr/bin/valgrind --xml=yes --xml-file=memory-report.xml \
	--leak-check=full --show-reachable=yes --show-leak-kinds=all \
	--errors-for-leak-kinds=all --error-exitcode=1 --track-origins=yes \
	/builds/iopsys/map-agent/src/mapagent -vvvvvv -d"
